# scripts/tenant-worker-reconciler.service
[Unit]
Description=PM Tenant Worker Reconciler
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/production-env
ExecStart=/opt/production-env/scripts/tenant-worker-reconciler.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pm-tenant-reconciler
