# scripts/tenant-worker-reconciler.timer
[Unit]
Description=PM Tenant Worker Reconciler Timer (30s)

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
